# scrapyd.service
[Unit]
Description = Scrapyd

[Service]
ExecStart = /usr/loca/bin/scrapyd
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
